# update-agent/debian/worldcoin-update-agent.timer
[Unit]
Description=Run update-agent after boot and every hour

[Timer]
OnBootSec=5min
OnUnitInactiveSec=5min

[Install]
WantedBy=timers.target
