# update-verifier/debian/orb-update-verifier.worldcoin-update-verifier.service
[Unit]
Description=Worldcoin Update Verifier Service.
StartLimitInterval=0
After=worldcoin-create-hardware-version.service
After=worldcoin-configure-can@can0.service

[Service]
Type=oneshot
User=root
ExecStart=/usr/local/bin/orb-update-verifier
Restart=on-failure
RestartSec=10s
SyslogIdentifier=worldcoin-update-verifier

[Install]
WantedBy=multi-user.target
