# ui/debian/orb-ui.worldcoin-ui.service
[Unit]
Description=Worldcoin UI Service
Documentation=https://github.com/worldcoin/orb-software
After=worldcoin-dbus.socket
Wants=worldcoin-dbus.socket
After=sound.target
Requires=sound.target

[Service]
Type=simple
User=orb-ui
Group=orb-ui
SupplementaryGroups=worldcoin-dbus
StateDirectory=orb-ui
Environment="HOME=/var/lib/orb-ui"
Environment="DBUS_SESSION_BUS_ADDRESS=unix:path=/tmp/worldcoin_bus_socket"
Environment="RUST_BACKTRACE=1"
Environment="PULSE_SERVER=unix:/run/user/1000/pulse/native"
SyslogIdentifier=worldcoin-ui
ExecStart=/usr/local/bin/orb-ui daemon
Restart=always

NoNewPrivileges=yes
ProtectHome=read-only
ProtectSystem=strict
ReadOnlyPaths=/
# PulseAudio creates its runtime dir in /tmp; D-Bus socket lives at /tmp/worldcoin_bus_socket
ReadWritePaths=/tmp
PrivateTmp=no
# Deny access to all devices by default; explicitly allow only what is needed
DevicePolicy=closed
# ALSA PCM playback (/proc/devices: 116 alsa)
DeviceAllow=char-alsa rw
# UART to MCU: Pearl uses ttyTHS0, Diamond uses ttyTHS1 (/proc/devices: 240 ttyTHS)
DeviceAllow=char-ttyTHS rw
# AF_UNIX: D-Bus session bus and PulseAudio socket
# AF_NETLINK: libasound hardware event monitoring
RestrictAddressFamilies=AF_UNIX AF_NETLINK
RestrictNamespaces=yes
RestrictSUIDSGID=yes
LockPersonality=yes
ProtectHostname=yes
ProtectKernelLogs=yes
ProtectKernelModules=yes
ProtectKernelTunables=yes
ProtectControlGroups=yes
ProtectProc=invisible
MemoryDenyWriteExecute=yes
SystemCallArchitectures=native
CapabilityBoundingSet=
SystemCallFilter=~@cpu-emulation @debug @module @mount @obsolete @raw-io @reboot @swap @privileged

[Install]
WantedBy=multi-user.target

